# packaging/debian/prerm
#!/bin/bash

systemctl disable agent_protocol_splitter.service micrortps_agent.service mavlink-router.service
systemctl disable mesh-init.service mesh.service mesh_pub.service
systemctl disable cloud_link.service mission_engine.service mission-data-recorder.service provisioning-agent.service
systemctl disable update_agent.service ota_update.service
systemctl disable mocap_pose.service
systemctl disable navigation.service control_interface.service octomap_server2.service bumper.service odometry2.service
systemctl disable depthai_ctrl.service depthai_gstreamer_node.service

exit 0
