# vm-comms/packaging/debian/postinst
#!/bin/bash

systemctl enable communication_link.service videonode.service
