# vm-comms/packaging/debian/postinst
#!/bin/bash

systemctl enable cloud_link.service

if [ -e /etc/systemd/system/mission_engine.service ]; then
    systemctl enable mission_engine.service
fi

exit 0
